# local/launchd/com.user.excel-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.excel-sync</string>

    <!-- После D1-centric pivot этот агент делает только бэкап D1 раз в день. -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/stepan/Desktop/excel/.venv/bin/python</string>
        <string>/Users/stepan/Desktop/excel/local/scripts/backup_d1.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/stepan/Desktop/excel</string>

    <key>RunAtLoad</key>
    <true/>

    <!-- Раз в день: 86400 секунд -->
    <key>StartInterval</key>
    <integer>86400</integer>

    <key>StandardOutPath</key>
    <string>/Users/stepan/Desktop/excel/local/logs/backup.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/stepan/Desktop/excel/local/logs/backup.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>
    <key>LowPriorityIO</key>
    <true/>
    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
